# dot_config/systemd/user/audio-keep-alive.service
[Unit]
Description=Keep audio alive

[Service]
ExecStart=/usr/bin/play -n -c2 synth sin gain

[Install]
WantedBy=default.target
